# broker/conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisord_bugel_broker.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord_bugel_broker.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord_bugel_broker.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///tmp/supervisord_bugel_broker.sock ; use a unix:// URL  for a unix socket

[program: broker]
process_name = broker
directory = /data/www/bugel_broker
command = /data/www/bugel_broker/bin/broker -c /data/www/bugel_broker/conf/broker.conf
numprocs=1
numprocs_start=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/bugel_broker/broker-console.log

[program: cross]
process_name = cross
directory = /data/www/bugel_broker
command = /data/www/bugel_broker/bin/cross_tcp -port=1881
numprocs=1
numprocs_start=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/bugel_broker/cross-console.log

